# systemd/genkai-tour-generate-world-list.service
[Unit]
Description=GenkaiTour Generate World List
After=network.target

[Service]
ExecStart=/usr/bin/make all
WorkingDirectory=/home/wacky/GenkaiTourWorldList
