# systemd/genkai-tour-partial-update-world-list.service
[Unit]
Description=GenkaiTour Generate World List
After=network.target

[Service]
ExecStart=/usr/bin/make partial-update
WorkingDirectory=/home/wacky/GenkaiTourWorldList
